# etc/systemd/system/co2mini.timer
[Unit]
Description=Timer to get SAMPLE values for zabbix agent

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timer.target
